# extension/edr/com.fleetdm.edr.notify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!--
	    Fleet EDR per-user notification surface (LaunchAgent).
	    Run by launchd in each user's session, this plist starts
	    the host app in `notify` mode so the system extension's
	    NotificationClient can post AUTH_EXEC-denial modals.

	    Install path (per user): ~/Library/LaunchAgents/com.fleetdm.edr.notify.plist
	    System-wide alternative: /Library/LaunchAgents/com.fleetdm.edr.notify.plist
	    The system-wide path requires a logout/login cycle on first install.
	-->
	<key>Label</key>
	<string>com.fleetdm.edr.notify</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Fleet EDR.app/Contents/MacOS/edr</string>
		<string>notify</string>
	</array>

	<!--
	    MachServices binds the XPC name the extension connects to.
	    launchd routes the first connection to this Mach service to
	    spawning the host app, so the notify surface is on-demand
	    rather than always-resident. The name matches the
	    blockNotificationServiceName constant baked into both
	    extension/NotificationClient.swift and edr/BlockNotification.swift.
	-->
	<key>MachServices</key>
	<dict>
		<key>FDG8Q7N4CC.com.fleetdm.edr.notifications</key>
		<true/>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<!--
	    KeepAlive=false: the agent re-spawns the process on demand
	    via the Mach service binding. A crashed host app stays down
	    until the next AUTH_EXEC denial, which is the desired
	    behaviour — we don't want a CPU-burning crashloop on the
	    user's desktop. The next denial spawns a fresh instance.
	-->
	<key>KeepAlive</key>
	<false/>

	<key>ProcessType</key>
	<string>Interactive</string>

	<!--
	    No StandardOutPath / StandardErrorPath: the host app routes
	    its output through os_log (Logger(subsystem: "com.fleetdm.edr"))
	    so anything worth keeping shows up in the unified log and is
	    queryable via `log show --predicate 'subsystem == "com.fleetdm.edr"'`.
	    A shared file path under /tmp would race across users on a
	    multi-user install (the first user's launchd creates the file
	    with permissions that lock out later sessions) — caught by
	    Gemini + Copilot on PR #157.

	    Packaging note: the dist/pkg path does NOT yet copy this
	    plist into /Library/LaunchAgents. The demo dry-run runs the
	    host app manually in a Terminal session
	    (`sudo /Applications/Fleet\ EDR.app/Contents/MacOS/edr notify`)
	    so the daemon's Mach lookup lands in the same bootstrap
	    namespace as the user's process. Production packaging hooks
	    are a follow-up tracked alongside the LaunchAgent installer
	    work.
	-->
</dict>
</plist>
